# metacoach.service
[Unit]
Description=MetaCoach FastAPI Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/metacoach
EnvironmentFile=/var/www/metacoach/.env
ExecStart=/var/www/chatbot/venv/bin/uvicorn api:app --host 127.0.0.1 --port 8101 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
